# packages/fml-autonomous-veh-deps/ros2_device_permissions.service
[Unit]
Description=Systemd service to ensure device permission are correctly set after Boot.

[Service]
Type=simple
ExecStart=/bin/bash /usr/bin/set_device_permissions.sh

[Install]
WantedBy=multi-user.target